# ouroboros.service.in
[Unit]
Description=Ouroboros IPC Resource Manager Daemon

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_SBINDIR@/irmd @CONFIGURE_STRING@
User=root

[Install]
WantedBy=multi-user.target
